# src/Misc/layoutbin/vsts.agent.service.template
[Unit]
Description={{Description}}
After=network.target

[Service]
ExecStart={{AgentRoot}}/runsvc.sh
Environment=LD_LIBRARY_PATH=/usr/local/lib:/usr/lib
User={{User}}
WorkingDirectory={{AgentRoot}}
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5min

[Install]
WantedBy=multi-user.target
